# cfg/dtdnssyncd.service
[Unit]
Description=dtdns ip updater
After=network.target

[Service]
ExecStart=/usr/sbin/dtdnssyncd
KillSignal=SIGINT
User=dtdnssync

[Install]
WantedBy=multi-user.target
